# commands/src/main/resources/supervisord.conf
[supervisord]
nodeamon=true

[program:postgres]
command=/etc/init.d/postgresql start


[program:tomee]
command=/usr/local/tomee/bin/catalina.sh run
